# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/usr/bin/xbindkeys
Name=xbindkeys the AutoHotkey on Linux
